# systemd/oceano-llama-swap.service
[Unit]
Description=Oceano llama-swap — chat model server on :8081 (one model resident at a time)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__OCEANO_USER__
# scripts/install.sh templates the __OCEANO_USER__ / __LLAMA_DIR__ tokens and the
# llama-swap binary path to the real user / LLAMA_DIR / installed binary.
# llama.cpp's shared libs sit next to llama-server; expose them so a relocated build
# (RUNPATH baked at build time) still resolves them — children inherit this.
Environment=LD_LIBRARY_PATH=__LLAMA_DIR__/build/bin
ExecStart=/usr/local/bin/llama-swap -config __LLAMA_DIR__/llama-swap.yaml -listen 127.0.0.1:8081 -watch-config
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
