# systemd/oceano.service
[Unit]
Description=Oceano — local AI agent engine (web UI + Telegram + scheduler + embeddings)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# scripts/install.sh templates the __OCEANO_USER__ / __OCEANO_ROOT__ tokens below
# to the real service user and install directory.
User=__OCEANO_USER__
WorkingDirectory=__OCEANO_ROOT__
EnvironmentFile=-__OCEANO_ROOT__/oceano.env
# Web UI binds localhost only. The agent can run shell commands — reach it over an
# SSH tunnel or Tailscale; do NOT bind 0.0.0.0 without adding auth.
Environment=OCEANO_WEB_HOST=127.0.0.1
Environment=OCEANO_WEB_PORT=8800
Environment=PYTHONUNBUFFERED=1
# Put the user's local bin on PATH so delegate_to_claude finds the `claude` CLI
# (the official installer drops it in ~/.local/bin, which systemd's default PATH omits).
# NB: use the literal home — %h on a SYSTEM unit resolves to /root, not User='s home.
Environment=PATH=__OCEANO_HOME__/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/snap/bin
# One process owns it all: web UI, Telegram bot, scheduler loop, and the embedding
# server (spawned as a child). Set OCEANO_EMBED_MANAGED=0 to run embed separately.
ExecStart=__OCEANO_ROOT__/venv/bin/python -m oceano.engine
Restart=on-failure
RestartSec=10

# This agent runs shell commands BOTH internet-facing (web + telegram) and with no
# human in the loop (scheduler). Even if the service user has passwordless sudo,
# neuter privilege escalation and make all of $HOME read-only except the dirs Oceano
# must write — this confines even run_shell (e.g. blocks writes to ~/.ssh/authorized_keys).
# Two HOME exceptions: <home>/.claude (Claude Code — used by delegate_to_claude / the
# skills + eval reviewers — writes its session/config state there) and the llama.cpp dir
# (Rivers downloads .gguf models into llama.cpp/models and edits llama.cpp/llama-swap.yaml;
# without write access there, UI model downloads/serve fail with EROFS). The leading '-'
# makes each optional, so a box without ~/.claude — or one pointed at a remote LLM with no
# local llama.cpp — still starts cleanly.
NoNewPrivileges=true
ProtectHome=read-only
ReadWritePaths=__OCEANO_ROOT__/workspace __OCEANO_ROOT__/data __OCEANO_ROOT__/skills -__OCEANO_LLAMA_DIR__ -__OCEANO_HOME__/.claude
PrivateTmp=true
# SIGTERM the engine first so it can shut the embed child down cleanly, then
# SIGKILL anything still in the cgroup after the grace period.
KillMode=mixed
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
